# rpi-print-server/dinoco-dashboard.service
# DINOCO B2B Print Dashboard — systemd unit (V.42)
#
# SANDBOX POLICY (systemd v256+ default-safe) — เหมือน dinoco-print.service
# ดู comments ใน dinoco-print.service สำหรับ sandbox rules ครบ
#
# ReadWritePaths เพิ่ม /var/log สำหรับ Flask logging (ถ้าต้องการ)
#
# install.sh จะ sed replace /home/dinocoth/rpi-print-server → actual install dir

[Unit]
Description=DINOCO B2B Print Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dinocoth
WorkingDirectory=/home/dinocoth/rpi-print-server
ExecStart=/home/dinocoth/rpi-print-server/venv/bin/python3 /home/dinocoth/rpi-print-server/dashboard.py --port 5555
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

# ── Sandbox — explicit เพื่อ immune ต่อ systemd default changes ──
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/home/dinocoth/rpi-print-server/logs /home/dinocoth/rpi-print-server/tmp /tmp /var/log
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
