# rpi-print-server/dinoco-dashboard.service
[Unit]
Description=DINOCO B2B Print Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dinocoth
WorkingDirectory=/home/dinocoth/DINOCO-SYSTEM-Front-End---Back-End/rpi-print-server
ExecStart=/home/dinocoth/DINOCO-SYSTEM-Front-End---Back-End/rpi-print-server/venv/bin/python3 /home/dinocoth/DINOCO-SYSTEM-Front-End---Back-End/rpi-print-server/dashboard.py --port 5555
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/home/dinocoth/DINOCO-SYSTEM-Front-End---Back-End/rpi-print-server/logs /home/dinocoth/DINOCO-SYSTEM-Front-End---Back-End/rpi-print-server/tmp /tmp /var/log
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
