# B2B Project/rpi-print-server/dinoco-print.service
[Unit]
Description=DINOCO B2B Print Client
After=network-online.target cups.service
Wants=network-online.target

[Service]
Type=simple
User=dinocoth
WorkingDirectory=/home/dinocoth/rpi-print-server
ExecStart=/home/dinocoth/rpi-print-server/venv/bin/python3 /home/dinocoth/rpi-print-server/print_client.py --daemon
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/home/dinocoth/rpi-print-server/logs /home/dinocoth/rpi-print-server/tmp /tmp
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
